# pi52/autostart/kiosk.desktop
[Desktop Entry]
Type=Application
Name=煮麵系統 Kiosk
Exec=/home/pi52/kiosk.sh
X-GNOME-Autostart-enabled=true
